# src/main/assembly/dist/install/dd-sword2.service
[Unit]
Description=Dd Sword2 Service

[Service]
ExecStart=/opt/dans.knaw.nl/dd-sword2/bin/dd-sword2 server /etc/opt/dans.knaw.nl/dd-sword2/config.yml
# Java returns 143 even if the SIGTERM was handled correctly.
SuccessExitStatus=143

User=dd-sword2
Group=dd-sword2

[Install]
WantedBy=multi-user.target
